# supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 150MB
logfile_backups=10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /go/src/github.com/primait/go-bin
nocleanup = true
strip_ansi = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[inet_http_server]
port = 127.0.0.1:9001
username = prima
password = prima

[program:dlx]
process_name=%(program_name)s_%(process_num)02d
command=dlx -c src/github.com/primait/go-bin/parameters_dev.yml
numprocs=1
autorestart=true
autostart=true
