# supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 150MB
logfile_backups=10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /go/src/github.com/primait/go-bin
nocleanup = true
strip_ansi = false

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:9001        ; (ip_address:port specifier, *:port for all iface)
username=prima              ; (default is no username (open server))
password=prima             ; (default is no password (open server))

[program:dlx]
process_name=%(program_name)s_%(process_num)02d
command=./dlx/dlx -c parameters_prod.yml
numprocs=1
autorestart=true
autostart=true

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
